# jarvis-chrome.service
[Unit]
Description=Jarvis shared Chromium/Chrome
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/launch-jarvis-chrome.sh
Restart=on-failure
RestartSec=2
TimeoutStopSec=10
KillMode=mixed

[Install]
WantedBy=default.target
